# launchd/com.corrhi.claude-autonomous.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.corrhi.claude-autonomous</string>
    <key>ProgramArguments</key>
    <array>
        <string>bash</string>
        <string>$HOME/.claude/scripts/autonomous-work.sh</string>
    </array>
    <!-- Poll every 30 minutes. Script checks credits and skips instantly
         if 5hr window is depleted. Runs burst work when credits available. -->
    <key>StartInterval</key>
    <integer>1800</integer>
    <key>WorkingDirectory</key>
    <string>$HOME</string>
    <key>StandardOutPath</key>
    <string>/tmp/claude-autonomous-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/claude-autonomous-stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>$HOME</string>
        <key>PATH</key>
        <string>$HOME/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>CORRHI_VAULT</key>
        <string>$VAULT</string>
    </dict>
</dict>
</plist>
